# scripts/com.user.icloud-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.icloud-health</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/wilber_chen/.local/share/chezmoi/scripts/icloud-sync.sh health 2&gt;&amp;1 || true</string>
    </array>
    <key>StartInterval</key>
    <integer>7200</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/wilber_chen/.local/log/dotfiles/icloud-health.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/wilber_chen/.local/log/dotfiles/icloud-health.err.log</string>
</dict>
</plist>
